# systemd-example/wallpaper.service
[Unit]
Description=Wallpaper service

[Service]
ExecStart=/usr/bin/bash /path-to-repo/process-wallpaper/updateWallpaper.sh
